# deployment/systemd.service
[Unit]
Description=Family First Scheduler API
Documentation=https://github.com/yourusername/scheduler-api
After=network.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/clearhive-scheduler-api
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=scheduler-api
Environment="NODE_ENV=production"
Environment="PORT=3000"

[Install]
WantedBy=multi-user.target
